# deploy/systemd/catyolo-backend.service
[Unit]
Description=CatYolo backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${CATYOLO_USER}
Group=${CATYOLO_GROUP}
WorkingDirectory=${REPO_INSTALL_DIR}/catyolo_ai_backend
EnvironmentFile=${CONFIG_DIR}/backend.env
Environment=PATH=${REPO_INSTALL_DIR}/catyolo_ai_backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=${REPO_INSTALL_DIR}/catyolo_ai_backend/.venv/bin/uvicorn main:app --host ${BACKEND_HOST} --port ${BACKEND_PORT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=catyolo-backend

[Install]
WantedBy=multi-user.target
